# mamba/etc/profile.d/02-env-python.sh
[ $(echo ${CONTAINER} | grep -c "python") ] || exit 0
