# conda/etc/services.d/jupyterlab/run
#!/usr/bin/with-contenv zsh

exec /opt/conda/bin/jupyter lab /workspace